# noninteractive/.crontab
#
# $Id$
#
MAILTO=kevin.lipscomb@kvrs.org
#
#mmhh DD MM d (mm=minute, hh=hour, DD=DayOfMonth, MM=MonthNum, d=DayOfWeekWhereZeroIsSunday)
#
# Daily
#
30 00  *  * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/OSCALERT/noninteractive/cad_activity_notification_agent.aspx
30 03  *  * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/OSCALERT/noninteractive/cad_activity_notification_agent.aspx
30 06  *  * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/OSCALERT/noninteractive/cad_activity_notification_agent.aspx
30 09  *  * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/OSCALERT/noninteractive/cad_activity_notification_agent.aspx
30 12  *  * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/OSCALERT/noninteractive/cad_activity_notification_agent.aspx
30 15  *  * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/OSCALERT/noninteractive/cad_activity_notification_agent.aspx
30 18  *  * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/OSCALERT/noninteractive/cad_activity_notification_agent.aspx
30 21  *  * * wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/OSCALERT/noninteractive/cad_activity_notification_agent.aspx
#
# Weekly
#
#
# Monthly
#
#
# Quarterly
#
#
# Yearly
#
#
# Other
#
@reboot sleep 240 && wget --output-document=/dev/null --no-check-certificate --timeout=14400 --tries=1 https://localhost/OSCALERT/noninteractive/cad_activity_notification_agent.aspx
